# systemd-alt-cron/lib/systemd/system/timer-xe.target
[Unit]
Description=Alternative cron timers for backups and defered service-start after boot

[Install]
WantedBy=timers.target
